# supervisord-serf.conf
[program:serf]
command=/start-serf.sh
numprocs=1
autostart=true
autorestart=true
